# systemd/mastodon-cron.timer
[Unit]
Description=Mastodon periodic cache cleaner

[Timer]
OnUnitActiveSec=1day
OnActiveSec=1day

[Install]
WantedBy=timers.target
